# nix-config/systemd/targets/multi-user.target
[Unit]
Description=Multi-User System
Documentation=man:systemd.special(7)
Requires=basic.target
After=basic.target
AllowIsolate=yes 