# nix-config/systemd/targets/rescue.target
[Unit]
Description=Rescue Mode
Requires=rescue.service
After=rescue.service 