# nix-config/systemd/targets/serial-getty.target
[Unit]
Description=Serial Console Getty
Documentation=man:systemd.special(7)
Requires=basic.target
After=basic.target
AllowIsolate=yes 