# nix-config/systemd/targets/sysinit.target
[Unit]
Description=System Initialization
Documentation=man:systemd.special(7)
Requires=local-fs.target
After=local-fs.target 